# supervisor/cron.conf
[program:cron]
command=go-cron -file=/etc/gocron/crontab
priority=100
directory=%(ENV_CRON_WORKDIR)s
stdout_logfile=syslog
stderr_logfile=syslog
