# supervisor/rsyslog.conf
[program:rsyslog]
command=rsyslogd -n
priority=10
stdout_logfile=syslog
stderr_logfile=syslog
